# deploy/novahome/authlyn-cert-sync.timer
[Unit]
Description=Nightly sync of the authlyn TLS cert from the Pi (FENRIR)

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
