# deploy/novahome/surrealdb.service
[Unit]
Description=SurrealDB (authlyn backing store, novahome)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=surrealdb
Group=surrealdb
WorkingDirectory=/data/surrealdb
# surrealkv store on the dedicated /data volume (vs the Pi's /var/lib/surrealdb).
# root/root is the v1 placeholder, same as the Pi; scoped creds land with auth.
ExecStart=/usr/local/bin/surreal start --user root --pass root --bind 127.0.0.1:8000 surrealkv:///data/surrealdb
Restart=on-failure
RestartSec=2s
TimeoutStartSec=15s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/data/surrealdb

[Install]
WantedBy=multi-user.target
