# deploy/authlyn-updater.service
[Unit]
Description=authlyn-interactive auto-deploy updater (poll GitHub Releases, swap, restart)
Documentation=https://github.com/DamieMoon/authlyn-interactive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=damien
Group=damien
ExecStart=/opt/authlyn/pi-updater.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=180s

# Minimal sandboxing — the script needs sudo and write access to
# /opt/authlyn via NOPASSWD. The timer surfaces persistent failures
# through journalctl rather than thrashing on retry.
NoNewPrivileges=false
ProtectSystem=false
