# deploy/authlyn.service
[Unit]
Description=authlyn-interactive (E2EE roleplay chat)
Documentation=https://github.com/DamieMoon/authlyn-interactive
Requires=surrealdb.service
After=surrealdb.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=authlyn
Group=authlyn
WorkingDirectory=/opt/authlyn
ExecStart=/opt/authlyn/authlyn-interactive
EnvironmentFile=/opt/authlyn/.env
Restart=on-failure
RestartSec=2s
TimeoutStartSec=15s
TimeoutStopSec=10s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/opt/authlyn/media

[Install]
WantedBy=multi-user.target
