# deploy/surrealdb.service
[Unit]
Description=SurrealDB (shared Pi infra)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=surrealdb
Group=surrealdb
WorkingDirectory=/var/lib/surrealdb
ExecStart=/usr/local/bin/surreal start --user root --pass root --bind 127.0.0.1:8000 surrealkv:///var/lib/surrealdb
Restart=on-failure
RestartSec=2s
TimeoutStartSec=15s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/var/lib/surrealdb

[Install]
WantedBy=multi-user.target
